# raspberry_pi/central_controller/central-controller.service
[Unit]
Description=Escape Room Central Controller
After=network.target mosquitto.service
Wants=network.target
Requires=mosquitto.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/Wizards/raspberry_pi/central_controller
ExecStart=/bin/bash /home/pi/Wizards/raspberry_pi/central_controller/start_central_controller.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/pi/Wizards
Environment=FLASK_APP=central_controller.py
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target